# deployments/server/netbox-tool.service
[Unit]
Description=NetBox Failover Tool Server
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=netbox-tool
Group=netbox-tool
WorkingDirectory=/opt/netbox-tool
EnvironmentFile=/etc/netbox-tool/netbox-tool.env
ExecStart=/opt/netbox-tool/bin/netbox-tool
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netbox-tool

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/netbox-tool /var/log/netbox-tool
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
